# hass-addon-sunsynk-multi/rootfs/etc/s6-overlay/s6-rc.d/sunsynk/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting the QSUSB64 add-on"
uv run python -m ha_addon_sunsynk_multi
