# hass-addon-esp/rootfs/etc/services.d/esp/run
#!/command/execlineb -P
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# Addon example: https://github.com/home-assistant/addons-example
# ==============================================================================

uv run python -m ha_addon_esp
